# xcontrol/postinst
systemctl enable docker
service docker start
exit 0;